# systemd/zig-fcps-api.service.template
[Unit]
Description=Scarlett FCP HTTP and MCP API daemon
After=network.target fcp-server@0.service
Requires=fcp-server@0.service

[Service]
Type=simple
User=zigfi
Group=audio
SupplementaryGroups=audio
ExecStart=@PREFIX@/bin/zig-fcps-api --listen 0.0.0.0 --port 8677 --card 0
Restart=on-failure
RestartSec=2
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
